# debian/postinst
#!/bin/bash
# postinst script for maxiofs
# This script runs after the package files are unpacked

set -e

case "$1" in
    configure)
        # Create maxiofs system user if it doesn't exist
        if ! getent passwd maxiofs > /dev/null 2>&1; then
            echo "Creating maxiofs system user..."
            adduser --system --group --home /var/lib/maxiofs \
                    --no-create-home --disabled-login \
                    --gecos "MaxIOFS Server" maxiofs
        fi

        # Set ownership of directories
        chown -R maxiofs:maxiofs /var/lib/maxiofs
        chown -R maxiofs:maxiofs /var/log/maxiofs

        # CRITICAL: Restore config.yaml from backup if it was deleted during upgrade
        # This handles the case where old packages included config.yaml and new ones don't
        if [ ! -f /etc/maxiofs/config.yaml ] && [ -f /etc/maxiofs/config.yaml.dpkg-backup ]; then
            echo "Restoring config.yaml from backup (preserving service credentials)..."
            cp -p /etc/maxiofs/config.yaml.dpkg-backup /etc/maxiofs/config.yaml
            rm -f /etc/maxiofs/config.yaml.dpkg-backup
            echo "✅ Configuration restored successfully!"
        elif [ ! -f /etc/maxiofs/config.yaml ]; then
            # Fresh install - create from example
            echo "Creating initial config.yaml from example..."
            cp /etc/maxiofs/config.example.yaml /etc/maxiofs/config.yaml
        else
            # Config exists and no backup needed
            echo "Preserving existing config.yaml (contains service credentials and secrets)"
            # Clean up backup if it exists
            rm -f /etc/maxiofs/config.yaml.dpkg-backup
        fi

        # Set ownership and permissions for config files
        if [ -f /etc/maxiofs/config.yaml ]; then
            chown maxiofs:maxiofs /etc/maxiofs/config.yaml
            chmod 0640 /etc/maxiofs/config.yaml
        fi

        if [ -f /etc/maxiofs/config.example.yaml ]; then
            chown maxiofs:maxiofs /etc/maxiofs/config.example.yaml
            chmod 0644 /etc/maxiofs/config.example.yaml
        fi

        # Set permissions
        chmod 0750 /var/lib/maxiofs
        chmod 0750 /var/log/maxiofs

        # Reload systemd
        if [ -d /run/systemd/system ]; then
            systemctl --system daemon-reload >/dev/null || true
        fi

        # Enable service (but don't start it automatically)
        if [ -d /run/systemd/system ]; then
            systemctl enable maxiofs.service >/dev/null || true

            # Check if this is a new installation or an upgrade
            if [ -z "$2" ]; then
                # New installation
                echo "MaxIOFS installed. Review /etc/maxiofs/config.yaml before starting."
                echo "Setup guide: /usr/share/doc/maxiofs/DEPLOYMENT.md"
                echo "After first start, download the encryption recovery bundle (see docs)."
            else
                # Upgrade — restart the service if it was running before the upgrade
                if [ -f /tmp/maxiofs-service-was-running ]; then
                    rm -f /tmp/maxiofs-service-was-running
                    echo "Restarting MaxIOFS service..."
                    systemctl start maxiofs.service >/dev/null || true
                fi

                echo "MaxIOFS upgraded. Your config.yaml was preserved."
                echo "Changes: /usr/share/doc/maxiofs/CHANGELOG.md"
            fi
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
