# debian/postinst
#!/bin/bash
# postinst script for maxiofs
# This script runs after the package files are unpacked

set -e

case "$1" in
    configure)
        # Create maxiofs system user if it doesn't exist
        if ! getent passwd maxiofs > /dev/null 2>&1; then
            echo "Creating maxiofs system user..."
            adduser --system --group --home /var/lib/maxiofs \
                    --no-create-home --disabled-login \
                    --gecos "MaxIOFS Server" maxiofs
        fi

        # Set ownership of directories
        chown -R maxiofs:maxiofs /var/lib/maxiofs
        chown -R maxiofs:maxiofs /var/log/maxiofs

        # Create config.yaml from example if it doesn't exist (NEVER overwrite existing config)
        if [ ! -f /etc/maxiofs/config.yaml ]; then
            echo "Creating initial config.yaml from example..."
            cp /etc/maxiofs/config.example.yaml /etc/maxiofs/config.yaml
        else
            echo "Preserving existing config.yaml (contains encryption keys)"
        fi

        # Set ownership and permissions for config files
        if [ -f /etc/maxiofs/config.yaml ]; then
            chown maxiofs:maxiofs /etc/maxiofs/config.yaml
            chmod 0640 /etc/maxiofs/config.yaml
        fi

        if [ -f /etc/maxiofs/config.example.yaml ]; then
            chown maxiofs:maxiofs /etc/maxiofs/config.example.yaml
            chmod 0644 /etc/maxiofs/config.example.yaml
        fi

        # Set permissions
        chmod 0750 /var/lib/maxiofs
        chmod 0750 /var/log/maxiofs

        # Reload systemd
        if [ -d /run/systemd/system ]; then
            systemctl --system daemon-reload >/dev/null || true
        fi

        # Enable service (but don't start it automatically)
        if [ -d /run/systemd/system ]; then
            systemctl enable maxiofs.service >/dev/null || true

            # Check if this is a new installation or an upgrade
            if [ -z "$2" ]; then
                # New installation
                echo ""
                echo "========================================"
                echo "MaxIOFS has been installed successfully!"
                echo "========================================"
                echo ""
                echo "Configuration file: /etc/maxiofs/config.yaml"
                echo "Example config: /etc/maxiofs/config.example.yaml"
                echo "Data directory: /var/lib/maxiofs"
                echo "Log directory: /var/log/maxiofs"
                echo ""
                echo "IMPORTANT: Before starting MaxIOFS, please:"
                echo "  1. Edit /etc/maxiofs/config.yaml with your settings"
                echo "  2. Ensure the data_dir is configured correctly"
                echo "  3. Review security settings (TLS, authentication)"
                echo ""
                echo "To start MaxIOFS:"
                echo "  sudo systemctl start maxiofs"
                echo ""
                echo "To check status:"
                echo "  sudo systemctl status maxiofs"
                echo ""
                echo "To view logs:"
                echo "  sudo journalctl -u maxiofs -f"
                echo ""
                echo "Web Console: http://localhost:8081"
                echo "S3 API: http://localhost:8080"
                echo ""
                echo "â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•"
                echo "âš ï¸  CRITICAL SECURITY WARNING âš ï¸"
                echo "â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•"
                echo "The config.yaml file contains your ENCRYPTION KEY."
                echo "If this key is lost, ALL encrypted data is PERMANENTLY LOST."
                echo ""
                echo "IMMEDIATELY create a backup:"
                echo "  sudo cp /etc/maxiofs/config.yaml /etc/maxiofs/config.yaml.backup"
                echo ""
                echo "Store the backup in a SECURE location (off-server recommended)"
                echo "â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•"
                echo ""
            else
                # Upgrade
                echo ""
                echo "========================================"
                echo "MaxIOFS has been upgraded successfully!"
                echo "========================================"
                echo ""
                echo "âœ… Your configuration has been PRESERVED: /etc/maxiofs/config.yaml"
                echo "   (Encryption keys and settings remain unchanged)"
                echo ""
                echo "ðŸ“„ Updated example config available at:"
                echo "   /etc/maxiofs/config.example.yaml"
                echo ""
                echo "ðŸ“ Data directory: /var/lib/maxiofs"
                echo ""
                echo "To restart MaxIOFS with the new version:"
                echo "  sudo systemctl restart maxiofs"
                echo ""
                echo "To check status:"
                echo "  sudo systemctl status maxiofs"
                echo ""
                echo "â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•"
                echo "âš ï¸  REMINDER: Backup your config.yaml regularly"
                echo "â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•"
                echo ""
            fi
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
